# systemd/agoogol-experiment@.timer
[Unit]
Description=Run Agoogol experiment daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
